# supervisor/calculate_k.conf
[group:calculate_k]
programs=tornado-8000,tornado-8001,tornado-8002,tornado-8003,tornado-8004,tornado-8005,tornado-8006,tornado-8007,tornado-8008,tornado-8009

[program:calculate-01]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-02]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-03]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-04]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-05]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-06]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-07]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-08]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-09]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-10]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-11]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-12]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-13]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-14]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-15]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-16]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-17]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-18]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-19]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-20]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-21]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-22]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-23]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:calculate-24]
command=/root/python/bin/python /home/ubuntu/package/kafka_calculate_k.py
directory=/root/python/bin
user=root
autorestart=true
stdout_logfile=/var/log/calculate_k/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/calculate_k/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info