# supervisor/tornado.conf
[group:tornadoes]
programs=tornado-8000,tornado-8001,tornado-8002,tornado-8003,tornado-8004,tornado-8005,tornado-8006,tornado-8007,tornado-8008,tornado-8009

[program:tornado-8000]
command=/home/ubuntu/env/env/bin/python /home/ubuntu/server.py --port=8000
directory=/home/ubuntu/env/env/bin
user=root
autorestart=true
stdout_logfile=/var/log/tornado/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/tornado/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:tornado-8001]
command=/home/ubuntu/env/env/bin/python /home/ubuntu/server.py --port=8001
directory=/home/ubuntu/env/env/bin
user=root
autorestart=true
stdout_logfile=/var/log/tornado/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/tornado/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:tornado-8002]
command=/home/ubuntu/env/env/bin/python /home/ubuntu/server.py --port=8002
directory=/home/ubuntu/env/env/bin
user=root
autorestart=true
stdout_logfile=/var/log/tornado/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/tornado/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:tornado-8003]
command=/home/ubuntu/env/env/bin/python /home/ubuntu/server.py --port=8003
directory=/home/ubuntu/env/env/bin
user=root
autorestart=true
stdout_logfile=/var/log/tornado/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/tornado/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info
[program:tornado-8004]
command=/home/ubuntu/env/env/bin/python /home/ubuntu/server.py --port=8004
directory=/home/ubuntu/env/env/bin
user=root
autorestart=true
stdout_logfile=/var/log/tornado/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/tornado/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:tornado-8005]
command=/home/ubuntu/env/env/bin/python /home/ubuntu/server.py --port=8005
directory=/home/ubuntu/env/env/bin
user=root
autorestart=true
stdout_logfile=/var/log/tornado/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/tornado/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:tornado-8006]
command=/home/ubuntu/env/env/bin/python /home/ubuntu/server.py --port=8006
directory=/home/ubuntu/env/env/bin
user=root
autorestart=true
stdout_logfile=/var/log/tornado/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/tornado/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:tornado-8007]
command=/home/ubuntu/env/env/bin/python /home/ubuntu/server.py --port=8007
directory=/home/ubuntu/env/env/bin
user=root
autorestart=true
stdout_logfile=/var/log/tornado/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/tornado/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:tornado-8008]
command=/home/ubuntu/env/env/bin/python /home/ubuntu/server.py --port=8008
directory=/home/ubuntu/env/env/bin
user=root
autorestart=true
stdout_logfile=/var/log/tornado/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/tornado/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info

[program:tornado-8009]
command=/home/ubuntu/env/env/bin/python /home/ubuntu/server.py --port=8009
directory=/home/ubuntu/env/env/bin
user=root
autorestart=true
stdout_logfile=/var/log/tornado/stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture_maxbytes=10MB
stderr_logfile=/var/log/tornado/stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture_maxbytes=10MB
loglevel=info